# deploy/systemd/uo_regulations_metadata_cleanup.timer.template
[Unit]
Description=Run UO_Regulations Metadata Cleanup on schedule

[Timer]
OnCalendar={{CLEANUP_ON_CALENDAR}}
Persistent=true

[Install]
WantedBy=timers.target
