# deploy/systemd/uo_regulations_backup.timer.template
[Unit]
Description=Run UO_Regulations Backup on schedule

[Timer]
OnCalendar={{BACKUP_ON_CALENDAR}}
Persistent=true

[Install]
WantedBy=timers.target
