# build/systemd/adoption-standard-update.timer
[Unit]
Description=Run Adoption Standard Update on schedule

[Timer]
OnCalendar=*-*-* 8:00:00
Persistent=true

[Install]
WantedBy=timers.target
